# dexedrine.service.template
[Unit]
Description=Dexedrine

[Service]
Type=simple
ExecStart=%DEXPATH%
KillSignal=SIGINT
